# scripts/homebox-vision.service
[Unit]
Description=Homebox Companion - AI-powered inventory assistant
After=network.target

[Service]
User=homebox
Group=homebox
WorkingDirectory=/home/homebox/homebox-app
EnvironmentFile=/etc/homebox.env

# Use uv to run with the managed virtual environment
ExecStart=/home/homebox/.local/bin/uv run python -m server.app

# Alternative if uv is installed elsewhere or you prefer direct venv:
# ExecStart=/home/homebox/homebox-app/.venv/bin/python -m server.app

Restart=always
RestartSec=5

# Security hardening (optional)
NoNewPrivileges=true
ProtectSystem=strict
ReadWritePaths=/home/homebox/homebox-app/logs

[Install]
WantedBy=multi-user.target






